# rc.d/dbmaster.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: dbmaster
# REQUIRE: LOGIN

. /etc/rc.subr

name="dbmaster"
rcvar="dbmaster_enable"

pidfile="@APP_RUNDIR@"/dbmaster.pid
command="@prefix@/sbin/${name}"
command_args=""
procname="@prefix@/sbin/${name}"
command_interpreter="@perl@"

load_rc_config ${name}

: ${dbmaster_enable:="NO"}

run_rc_command "$1"
#EOF
