# rc.d/pgmaster.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: pgmaster
# REQUIRE: LOGIN

. /etc/rc.subr

name="pgmaster"

rcvar="pgmaster_enable"

pidfile="@APP_RUNDIR@"/pgmaster.pid
command="@prefix@/sbin/${name}"
procname="@prefix@/sbin/${name}"
command_interpreter="@PERL@"


load_rc_config ${name}

: ${pgmaster_enable:="NO"}

run_rc_command "$1"
#EOF
