# rc.d/pgstore.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: pgstore
# REQUIRE: LOGIN

. /etc/rc.subr

name="pgstore"

rcvar="pgstore_enable"

pidfile="@APP_RUNDIR@"/pgstore.pid
command="@prefix@/sbin/${name}"
procname="@prefix@/sbin/${name}"
command_interpreter="@PERL@"


load_rc_config ${name}

: ${pgstore_enable:="NO"}

run_rc_command "$1"
#EOF
