# dbmaster.service.in
[Unit]
Description=dbmaster

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/dbmaster.pid
ExecStart=@prefix@/sbin/dbmaster
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

