# pgagent.service.in
[Unit]
Description=PGagent

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/pgagent.pid
ExecStart=@APP_EXECDIR@/pgagent
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

