# pgmaster.service.in
[Unit]
Description=PGmaster

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/pgmaster.pid
ExecStart=@APP_EXECDIR@/pgmaster
ExecReload=/bin/kill -HUP $MAINPID


[Install]
WantedBy=multi-user.target

